# com.akitomo.st_channel.rss.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.akitomo.st_channel.rss</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/akitomo/src/vault/10_Projects/ST_channnel/99_System/.venv/bin/python3</string>
        <string>/Users/akitomo/src/vault/10_Projects/ST_channnel/99_System/scripts/fetch_rss_links.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/com.akitomo.st_channel.rss.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.akitomo.st_channel.rss.err</string>
</dict>
</plist>
